# STM32MP1_SSK_Pkg_Rel/SSK_Suit_Configuration/Tpm_Measured_Boot/rc.local
#!/bin/bash
/etc/init.d/Measured_boot.sh

